# board/common/rootfs-overlay/etc/systemd/system/openhab.service.d/overlay.conf
[Service]
EnvironmentFile=-/boot/openhab.txt
