# board/common/rootfs-overlay/etc/systemd/system/sshd.service.d/override.conf
[Unit]
Wants=ssh-host-key@rsa.service ssh-host-key@ecdsa.service ssh-host-key@ed25519.service

[Service]
ExecStartPre=
StateDirectory=ssh
